# examples/backup-cleaner.service
[Unit]
Description=Remove old backups

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cleaner.py /media/hdd/snapshots
