# examples/backup-btrfs.service
[Unit]
Description=Backup filesystem

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup-btrfs.sh
